# tools/nightly/com.philtullai.av-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template launchd plist for the nightly test-fix job. Install with:

    cp tools/nightly/com.philtullai.av-nightly.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.philtullai.av-nightly.plist

  Fires at 3:00 AM local time. Logs to ~/av_runtime/logs/launchd/nightly.{out,err}
  and ~/av_runtime/logs/nightly/<date>.log.

  Uninstall:

    launchctl unload ~/Library/LaunchAgents/com.philtullai.av-nightly.plist
    rm ~/Library/LaunchAgents/com.philtullai.av-nightly.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.philtullai.av-nightly</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/philtullai/av/tools/nightly/run.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/philtullai/av</string>
	<key>StandardOutPath</key>
	<string>/Users/philtullai/av_runtime/logs/launchd/nightly.out</string>
	<key>StandardErrorPath</key>
	<string>/Users/philtullai/av_runtime/logs/launchd/nightly.err</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
